# services/rc-local.service
[Unit]
Description=Run rc.local

[Service]
ExecStart=/etc/rc.local
Type=oneshot

[Install]
WantedBy=multi-user.target
